# IML.ScannerProxyDaemon/systemd-units/scanner-proxy.path
[Unit]
Description=IML Scanner Proxy Path
PartOf=device-scanner.socket
BindsTo=device-scanner.service
After=device-scanner.service

[Path]
DirectoryNotEmpty=/etc/iml

[Install]
WantedBy=device-scanner.socket
